# etc/systemd/event-driven-system.service
# Event-Driven Repository Understanding System
# NOTE: This file contains hard-coded paths that must be updated for your environment.
# Before installing, update the following:
# 1. User: Replace 'runner' with your actual username
# 2. WorkingDirectory: Update to your repository's absolute path
# 3. ExecStart: Update to match your repository's absolute path
# 4. StandardOutput/StandardError: Update log paths to match your repository

[Unit]
Description=Event Driven Repository Understanding System
After=network.target

[Service]
Type=simple
# Replace with your actual user
User=runner
# Replace with the absolute path to your repository
WorkingDirectory=/home/runner/work/machine-native-ops/machine-native-ops/workspace/tools/repository-understanding
ExecStart=/usr/bin/python3 /home/runner/work/machine-native-ops/machine-native-ops/workspace/tools/repository-understanding/event_driven_system.py
Restart=always
RestartSec=10
StandardOutput=append:/home/runner/work/machine-native-ops/machine-native-ops/workspace/tools/repository-understanding/logs/system.log
StandardError=append:/home/runner/work/machine-native-ops/machine-native-ops/workspace/tools/repository-understanding/logs/error.log

[Install]
WantedBy=multi-user.target
